# contrib/etc/systemd/system/nginx.service.d/override.conf
[Service]
PIDFile=/run/nginx/nginx.pid
ExecStartPre=
ExecStartPre=/usr/bin/rm -f /run/nginx/nginx.pid
ExecStartPre=/usr/sbin/nginx -t
